# modules/bootstrap/templates/systemd/zfs-import-home.service.erb
[Unit]
Description=Load <%= @rpool %>/home encryption keys
Before=systemd-user-sessions.service
Before=zfs-mount.service
After=zfs-import.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c 'systemd-ask-password "Encrypted ZFS password for <%= @zfs_ds_home %>" --no-tty | zfs load-key <%= @zfs_ds_home %>'

[Install]
WantedBy=zfs-mount.service